# config/systemd/user/zeitgeist-daemon.service
[Unit]
Description=Zeitgeist Daemon

[Service]
ExecStart=/usr/bin/zeitgeist-daemon

[Install]
WantedBy=default.target
